# crontab/backup_keepass.sh
#!/bin/bash

dt=$(date '+%d_%m_%Y__%H:%M:%S')
cp /Users/marcus/Documents/keypass/marcus.kdbx /Users/marcus/Documents/keypass/backups/$dt.backup.kdbx
